# notion_pdf/scripts/macos/com.notion_pdf.cloudflare.quick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.notion_pdf.cloudflare.quick</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Volumes/mac_ssd/external_disk/Git/dev/projects/Tilon/notion_pdf/scripts/macos/start_cloudflare_quick_tunnel.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Volumes/mac_ssd/external_disk/Git/dev/projects/Tilon/notion_pdf</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/notion_pdf_cloudflared.quick.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/notion_pdf_cloudflared.quick.launchd.err.log</string>
</dict>
</plist>
